# deployment/folio-migration-web.service
[Unit]
Description=FOLIO Migration Web Application
After=network.target
Wants=network-online.target

[Service]
Type=exec
User=folio
Group=folio
WorkingDirectory=/opt/folio_migration_web

# Environment
Environment="PATH=/opt/folio_migration_web/.venv/bin:/usr/local/bin:/usr/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/opt/folio_migration_web/.env

# Start command
ExecStart=/opt/folio_migration_web/.venv/bin/uvicorn \
    folio_migration_web.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 4 \
    --log-level info

# Restart policy
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/folio_migration_web

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
