# deployment/folio-migration-web.service
[Unit]
Description=FOLIO Migration Web Application
After=network.target
Wants=network-online.target

[Service]
Type=exec
User=root
Group=root
WorkingDirectory=/folio/folio_migration_web

# Environment
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/folio/folio_migration_web/.env

# Start command using bash to properly resolve venv
ExecStart=/bin/bash -c 'cd /folio/folio_migration_web && source .venv/bin/activate && exec uvicorn folio_migration_web.main:app --host 127.0.0.1 --port 8000 --workers 4 --log-level info'

# Restart policy
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Security hardening (relaxed for uv python in /root)
NoNewPrivileges=yes
PrivateTmp=yes

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
